# service_files/xyz.openbmc_project.License.LicenseControl.service
[Unit]
Description=License Control service
After=systemd-timesyncd.service
StartLimitBurst=5
StartLimitInterval=5min
StartLimitAction=reboot

[Service]
ExecStart=/usr/bin/env licensecontrol
SyslogIdentifier=license-control
Restart=always
RestartSec=60
Type=dbus
BusName=xyz.openbmc_project.License

[Install]
WantedBy=multi-user.target

